# dot_config/systemd/user/battery-notify.service
[Unit]
Description=Battery low notification daemon
After=graphical-session.target

[Service]
ExecStart=%h/.local/bin/battery-notify
Restart=on-failure

[Install]
WantedBy=graphical-session.target
